# openntpd/openntpd.service
[Unit]
Description=OpenNTP daemon
After=network.target
 
[Service]
Type=forking
ExecStart=/usr/sbin/openntpd

[Install]
WantedBy=multi-user.target
